# net-misc/ddns_updater/client/ddns
#!/sbin/openrc-run
# (c) joe brendler 2026-7689
# this is the openrc service companion to /usr/sbin/ddns-daemon and /usr/sbin/ddns-update
# (this script should be installed as /etc/init.d/ddns

description="IPv6 SLAAC DDNS watcher"

supervisor="supervise-daemon"

command="/usr/sbin/ddns-daemon"
#command_background="yes"
pidfile="/run/ddns-daemon.pid"

depend() {
    use net
}

