# net-misc/ddns_updater/client/ddns
#!/sbin/openrc-run
# (c) joe brendler 2026-7689
# this is the openrc service companion to /usr/sbin/ddns-daemon and /usr/sbin/ddns-update
# (this script should be installed as /etc/init.d/ddns

description="IPv6 SLAAC DDNS daemon supervisor"
#command_background="yes"
# NOT using command_background - leaves daemon child process running when stopped
# using supervise-daemon instead, because ithandles backgrounding and pidfile automatically
supervisor="supervise-daemon"
command="/usr/sbin/ddns-daemon"
pidfile="/run/ddns-daemon.pid"

# tune supervise-daemon behavior:
# adding respawn_delay and respawn_max to the init script provides a safety net
#   for the daemon's lifecycle. supervise-daemon monitors the pid of ddns-daemon, and
#   if the ddns-daemon crashes, supervise-daemon will try to restart it, making it more
#   resilient; adding respawn_delay and _max provices rapid-fire crash throttling
#   to prevent it from sucking up all available CPU just crashing repeatedly
respawn_delay=5
respawn_max=10


depend() {
    #need net  ## no - ddns facilitates net functionality,
    # so doing "need net" cause "blocage"
    # do "use net" instead
    use net
}

