# net-misc/ddns_updater/client/ddns
#!/sbin/openrc-run
# (c) joe brendler 2026-7689
# this is the openrc service companion to /usr/sbin/ddns-daemon and /usr/sbin/ddns-update
# (this script should be installed as /etc/init.d/ddns

description="IPv6 SLAAC DDNS daemon supervisor"
#command_background="yes"
# NOT using command_background - leaves daemon child process running when stopped
# using supervise-daemon instead, because ithandles backgrounding and pidfile automatically
supervisor="supervise-daemon"
command="/usr/sbin/ddns-daemon"
pidfile="/run/ddns-daemon.pid"

# tune supervise-daemon behavior:
# adding respawn_delay and respawn_max to the init script provides a safety net
#   for the daemon's lifecycle. supervise-daemon monitors the pid of ddns-daemon, and
#   if the ddns-daemon crashes, supervise-daemon will try to restart it, making it more
#   resilient; adding respawn_delay and _max provices rapid-fire crash throttling
#   to prevent it from sucking up all available CPU just crashing repeatedly
# For a background network daemon, I want to be resilient but eventually stop if it's broken
#   (throttle crash-restart if there is a problem, but be resilient via a wide respawn period)
#
# respawn_delay - (cool down time) # of seconds supervisor waits after crash before
#   trying to start daemon again
respawn_delay=10
# respawn_max - (strike out count) max # of times daemon can die w/in respawn period
#   before supervisor marks it "crashed"
respawn_max=5
# respawn_period - (window of observation) rolling window (sec) during which
#   supervisor counts how many times daemon has crashed
respawn_period=120
#
# Safety Net: If a catastrophic error occurs (e.g., a missing library or syntax error),
#   the service will stop and stay in a "crashed" state after about 50 seconds
#   (5 attempts), rather than looping infinitely and flooding logs
# Resilience: This is robust enough to survive transient network blips or temporary
#   kernel resource exhaustion that might cause ip monitor (daemon logic) to trip

depend() {
    #need net  ## no - ddns facilitates net functionality,
    # so doing "need net" cause "blocage"
    # do "use net" instead
    use net
}

