# dev-sbc/sbc-status-leds/sbc-status-leds.crontab
*/3 * * * *   root   /usr/sbin/sbc-status-leds

