# net-misc/ddns/client/ddns-parts-wsl/ddns-trigger-wsl.crontab
# trigger ever 5 minutes and let ddns-update determine via cache if update is needed
*/5 * * * * root /usr/sbin/ddns-trigger-wsl wsl-cron
